# debian/lightfield.postinst.in
#!/bin/sh
# postinst script for lightfield
#
# see: dh_installdeb(1)
#
#set -e
#
# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

RELEASE_TRAIN=@@RELEASE_TRAIN@@
BUILDTYPE=@@BUILDTYPE@@
VERSION=@@VERSION@@
ARCHITECTURE=@@ARCHITECTURE@@
if [ "${RELEASE_TRAIN}" = "base" ]
then
    SUFFIX=${BUILDTYPE}
else
    SUFFIX=${RELEASE_TRAIN}-${BUILDTYPE}
fi

case "$1" in
    configure)
        perl -lpi -e 's/^(?!##LF## )/##LF## /;' /etc/apt/sources.list /etc/apt/sources.list.d/* 2>/dev/null

        echo "deb file:/var/lib/lightfield/software-updates/lightfield-${SUFFIX}_${VERSION}_${ARCHITECTURE} ./" > /etc/apt/sources.list.d/volumetric-lightfield.list
        chown lumen:lumen /etc/apt/sources.list.d/volumetric-lightfield.list
        chmod 644 /etc/apt/sources.list.d/volumetric-lightfield.list

        udevadm control --reload
        udevadm trigger /dev/ttyACM* /dev/ttyUSB*

        systemctl daemon-reload
        systemctl set-default multi-user.target
        systemctl enable --no-ask-password --now set-projector-power.service clean-up-mount-points.service
        systemctl enable --no-ask-password getty@tty1.service
        systemctl daemon-reload

        if [ "${RELEASE_TRAIN}" = "dlp4710" ] || [ "${RELEASE_TRAIN}" = "xdlp4710-20um" ]
        then
            set-projector-power --first-time
        fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument '$1'" 1>&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
